# luci-theme-m3e/root/etc/uci-defaults/30_luci-theme-m3e
#!/bin/sh

# M3E Theme UCI configuration
# Use option names without '-' for better UCI compatibility.
uci set luci.themes.M3E=/luci-static/m3e
uci set luci.themes.M3EBlue=/luci-static/m3e-blue
uci set luci.themes.M3EGreen=/luci-static/m3e-green
uci set luci.themes.M3ERed=/luci-static/m3e-red

# Clean up legacy keys from previous naming experiments.
uci -q delete luci.themes.M3E-Blue
uci -q delete luci.themes.M3E-Green
uci -q delete luci.themes.M3E-Red

if [ "$PKG_UPGRADE" != 1 ]; then
	uci set luci.main.mediaurlbase=/luci-static/m3e
fi

uci commit luci

exit 0